# etc/profile.d/nts-00-local-path.sh
# sourced by /etc/profile
case ":$PATH:" in *":/opt/nts/bin:"*) ;; *) PATH="$PATH:/opt/nts/bin";; esac
export PATH
